# deploy/systemd/micu-market-notification-emails.service
[Unit]
Description=Micu Market notification email dispatcher

[Service]
Type=oneshot
User=micu
Group=www-data
WorkingDirectory=/home/micu/Micu_market
EnvironmentFile=/home/micu/Micu_market/.env
Environment=DJANGO_SETTINGS_MODULE=Micu_market.settings_production
ExecStart=/home/micu/Micu_market/venv/bin/python /home/micu/Micu_market/manage.py send_notification_emails --limit 200 --settings=Micu_market.settings_production
