# deploy/systemd/micu-market-media-cleanup.timer
[Unit]
Description=Run Micu Market orphan media cleanup monthly

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
