# deploy/systemd/micu-market.service
[Unit]
Description=Micu Market Gunicorn service
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=micu
Group=www-data
UMask=0007
WorkingDirectory=/home/micu/Micu_market
EnvironmentFile=/home/micu/Micu_market/.env
Environment=DJANGO_SETTINGS_MODULE=Micu_market.settings_production
ExecStart=/home/micu/Micu_market/venv/bin/gunicorn -c /home/micu/Micu_market/gunicorn.conf.py Micu_market.asgi:application
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/micu/Micu_market/media /home/micu/Micu_market/staticfiles /var/log/micu_market

[Install]
WantedBy=multi-user.target
